# etc/supervisor.d/worker.ini
[program:symfony-worker]
process_name=%(program_name)s_%(process_num)02d
command=/home/siplah/app/bin/console --env=prod --no-debug enqueue:consume
autostart=true
autorestart=true
priority=15
numprocs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
